# .config/systemd/user/filen-sync.service
[Unit]
Description=Filen Obsidian Sync Daemon
After=network-online.target

[Service]
Type=simple
# Ensure this points to where your filen binary actually is
ExecStart=/usr/bin/filen-cli sync
# If it crashes, wait 10s and restart
Restart=always
RestartSec=10s

[Install]
WantedBy=default.target
